# packages/journal-to-fedora-messaging/journal_to_fedora_messaging-1.0.0.tar.gz/journal_to_fedora_messaging-1.0.0/journal-to-fedora-messaging.service
[Unit]
Description=Journal to Fedora Messaging
Documentation=https://github.com/fedora-infra/journal-to-fedora-messaging
Wants=network.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/journal-to-fedora-messaging
Restart=on-failure
RestartSec=60
User=journal2fedmsg
Group=journal2fedmsg

[Install]
WantedBy=multi-user.target
