# packages/FabOMatic/fabomatic-0.7.2.tar.gz/fabomatic-0.7.2/doc/config/systemd/user/default.target.wants/FabOMatic.service
[Unit]
Description=FabOMatic backend service
After=multi-user.target
After=network.target

[Service]
Type=exec
Delegate=true
KillMode=process
ExecStartPre=pip install FabOMatic --upgrade
ExecStart=python -m FabOMatic

[Install]
WantedBy=default.target	
